# cron.hourly/script
#!/bin/bash

sudo python /home/dietpi/InfoClimat.py
