# etc/cron.hourly/script
#!/bin/bash
python /home/dietpi/MeteoBzh.py
#python /home/dietpi/InfoClimatPast.py
/home/dietpi/MosquittoSub.sh
